# etc/supervisord.conf
[program:askgraph]
command=/home/ubuntu/askgraph.com/main.py
