# cron/askgraph.cron
0 0 * * * /home/ubuntu/askgraph.com/cron/backup_database.py
