# cron/askgraph.cron
0 0 * * * /home/ubuntu/askgraph.com/cron/backup_database.py
0 0 * * * /home/ubuntu/askgraph.com/cron/backup_sources.py
1 * * * * /home/ubuntu/askgraph.com/cron/send_mail.py
3 3 * * * /home/ubuntu/askgraph.com/cron/sensors/stock/poll_stock_prices.py
